# php-laraval-app/docker/php-fpm/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/www/html/storage/logs/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
pidfile=/var/run/supervisord.pid

[program:queue-schedules]
process_name=%(program_name)s-%(process_num)02d
command=php artisan queue:work --queue=schedules --memory=128
autostart=true
autorestart=true
numprocs=1
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/queue-schedules.log
stopwaitsecs=60

[program:queue-notifications]
process_name=%(program_name)s-%(process_num)02d
command=php artisan queue:work --queue=notifications --memory=128 --timeout=120
autostart=true
autorestart=true
numprocs=1
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/queue-notifications.log
stopwaitsecs=60

[program:queue-default]
process_name=%(program_name)s-%(process_num)02d
command=php artisan queue:work --queue=default --memory=128 --timeout=120
autostart=true
autorestart=true
numprocs=1
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/queue-default.log
stopwaitsecs=60